# packages/static/kairos-overlay-files/files/etc/systemd/system/k0scontroller.service.d/override.conf
[Unit]
ConditionFileIsExecutable=
ConditionFileIsExecutable=/usr/bin/k0s

[Service]
ExecStart=
ExecStart=/usr/bin/k0s controller --debug=true --single=true --verbose=true
