# packages/system/init-svc/openrc/cos-setup/files/etc/init.d/cos-setup-boot
#!/sbin/openrc-run

depend() {
  provide cos-setup-boot
}

start() {
  kairos-agent run-stage boot
  eend 0
}