# packages/system/init-svc/openrc/cos-setup/files/etc/init.d/cos-setup-network
#!/sbin/openrc-run

depend() {
  after net
  provide cos-setup-network
}

start() {
  kairos-agent run-stage network
  eend 0
}