# packages/utils/keylime-agent/var-lib-keylime-secure.mount
[Unit]
Description=Kelime configuration filesystem
Before=keylime-agent.service

[Mount]
What=tmpfs
Where=/var/lib/keylime/secure
Type=tmpfs
Options=mode=0700,size=1m,uid=tss,gid=tss

[Install]
WantedBy=multi-user.target