# packages/static/kairos-overlay-files/files/etc/systemd/system/tpm2-abrmd.service.d/override.conf
# Don't wait for 90 seconds for tpm devices when they don't exist
[Service]
TimeoutSec=5

[Unit]
JobTimeoutSec=5
