# packages/system/init-svc/systemd/cos-setup-boot.service
[Unit]
Description=cOS system configuration
Before=getty.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/kairos-agent run-stage boot

[Install]
WantedBy=multi-user.target
