# packages/system/init-svc/systemd/cos-setup-network.service
[Unit]
Description=cOS setup after network
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/kairos-agent run-stage network

[Install]
WantedBy=multi-user.target
